# configs/etc/profile.d/update-path.sh
export PATH="/opt/bin:/vagrant/bin:$PATH"
